# data-collector-service/crontab

## Damn this crontab

*/5 * * * * /usr/local/bin/python /app/app.py 
